# rc.d/sqm.in
#!/bin/sh
#
# $Id$
#
# PROVIDE: sqm
# REQUIRE: LOGIN

. /etc/rc.subr

name="sqm"
rcvar="sqm_enable"

pidfile="@APP_RUNDIR@"/sqm.pid
command="@prefix@/sbin/${name}"
command_args=""
procname="@prefix@/sbin/${name}"
command_interpreter="@perl@"

load_rc_config ${name}

: ${sqm_enable:="NO"}

run_rc_command "$1"
#EOF
